# docker/etc/service/cpulimit/run
#!/bin/sh
exec cpulimit --limit=80 --exe=/home/bitcoin/bitcoin-0.10.0/bin/bitcoind
